# packaging/systemd/aegis.service
# Aegis Control Plane — Systemd Service
#
# Run the Aegis control plane API server with embedded UI.
#
# Installation:
#   sudo cp aegis.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now aegis
#
# Status:
#   sudo systemctl status aegis
#   journalctl -u aegis -f
#
# Prerequisites:
#   - Aegis binary at /usr/local/bin/aegis
#   - Run 'aegis bootstrap --production' first (creates config + DB)
#   - Caddy installed (optional, for port 80 public access)

[Unit]
Description=Aegis Infrastructure Gateway Control Plane
Documentation=https://github.com/aegis/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/aegis

# Environment file (optional overrides)
EnvironmentFile=-/etc/aegis/aegis.conf

# Pre-flight: ensure data directory and DB exist
ExecStartPre=/usr/bin/test -d /var/lib/aegis
ExecStartPre=/usr/bin/test -d /etc/aegis

# Start Aegis — config auto-discovered from:
#   1. /etc/aegis/config.yaml (production)
#   2. $HOME/.aegis/config.yaml
#   3. DefaultConfig() fallback
ExecStart=/usr/local/bin/aegis serve

# Restart behavior
Restart=on-failure
RestartSec=5

# Startup timing
TimeoutStartSec=15
TimeoutStopSec=10

# Security hardening
# Note: ProtectSystem=full mounts /usr, /boot, /etc read-only.
# Middleware install via UI won't work through systemd — use SSH for apt-get.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/aegis /etc/aegis /etc/caddy /etc/haproxy
ReadOnlyPaths=/usr/local/bin/aegis

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
